# docker/etc/service/rss-server/run
#!/bin/sh
cp -n /var/www/html/rss.db /data/rss.db || true