# docker/etc/service/rss-server/run
#!/bin/sh
if [ ! -f /data/rss.db ]; then
    cp /var/www/html/rss.db /data/rss.db
fi # TODO: check if database format changed and upgrade it